# xdns/files/xdns.init
#!/bin/sh /etc/rc.common

START=99
SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

start_xdns() {

        ADDR="10.89.1.200:8083"
        GATEWAY="10.89.1.1"
        VPN="hkvpn-lt2p"
        service_start /usr/bin/xdns -l $ADDR -f $GATEWAY -vpn $VPN > /dev/null 2>&1
        return $?
}

start() {
        start_xdns
}

stop() {
        service_stop /usr/bin/xdns
}



