# xdns/files/xdns.init
#!/bin/sh /etc/rc.common

START=99
SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

start_xdns() {

        ADDR="192.168.1.2:8083"
        service_start /usr/bin/xdns -l $ADDR > /dev/null 2>&1
        return $?
}

start() {
        start_xdns
}

stop() {
        service_stop /usr/bin/xdns
}



